# punk-records-mcp/com.punkrecords.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.punkrecords.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/macbook/.nvm/versions/node/v22.19.0/bin/npm</string>
        <string>run</string>
        <string>backup</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macbook/punk-record/punk-records-mcp</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/macbook/punk-records/.brain/backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macbook/punk-records/.brain/backup.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
